# dev-libs/libcgroup/files/cgrules.socket
[Unit]
Description=Control Group Rules Daemon Listen Socket

[Socket]
ListenStream=/run/cgred.socket

[Install]
WantedBy=sockets.target
