# net-analyzer/prometheus-alertmanager/files/prometheus-alertmanager-initd
#!/sbin/openrc-run
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

description="Prometheus Alertmanager"

pidfile="/run/${SVCNAME}.pid"
command=/usr/bin/prometheus-alertmanager
command_args="${PROMETHEUS_ALERTMANAGER_ARGS}"
command_background="true"

start_stop_daemon_args="-u ${PROMETHEUS_ALERTMANAGER_USER} -1 ${PROMETHEUS_ALERTMANAGER_LOG_FILE} -2 ${PROMETHEUS_ALERTMANAGER_LOG_FILE} -d ${PROMETHEUS_ALERTMANAGER_DATA_DIR}"

depend() {
	need net
}
